# templates/opencode@.service
[Unit]
Description=OpenCode Server - Preset %i
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60s
StartLimitBurst=5

[Service]
Type=simple
User=fine
Group=fine
WorkingDirectory="/home/fine"

EnvironmentFile=/etc/opencode/opencode-%i.env

# OPENCODE_BIN_PATH is dynamically injected by install.sh
ExecStart="OPENCODE_BIN_PATH" serve --port ${PORT} --hostname 127.0.0.1

Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
